# site/profile/files/applications/glances/glances-web.service
[Unit]
Description=Glances Web Server

[Service]
Type=simple
User=root
ExecStart=/usr/local/bin/glances -w
Restart=on-failure