# site/profile/files/oit/keep-mouse-moving.service
[Unit]
Description=Keep Mousing Moving

[Service]
Type=simple
ExecStart=/usr/local/bin/keep-mouse-moving
Restart=on-failure
